# SSH/etc/systemd/system/update-knockd.timer
[Unit]
Description=Run knockd port update daily

[Timer]
OnCalendar=00:00
Persistent=true

[Install]
WantedBy=timers.target
